# dotfiles_linux/.config/autostart/CAPS to ESCAPE remap.desktop
[Desktop Entry]
Exec=setxkbmap -option caps:escape
Name=CAPS_to_ESCAPE
Type=Application
Version=1.0
