# dotfiles_linux/.config/autostart/Redshift.desktop
[Desktop Entry]
Exec=redshift
Name=Redshift
Type=Application
Version=1.0
